# tests/host/with_key/docker-swarm-init.service
[Unit]
Description=Initialize Docker Swarm on first boot
After=docker.service
Wants=docker.service
ConditionPathExists=!/etc/docker/swarm_initialized

[Service]
Type=oneshot
ExecStart=/bin/bash -c "docker swarm init && touch /etc/docker/swarm_initialized"
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
